# outlayer-landing.service
[Unit]
Description=OutLayer Landing Page
After=network.target

[Service]
Type=simple
User=nextjs-user
WorkingDirectory=/home/nextjs-user/code/outlayer.ai
ExecStart=/usr/local/bin/node server.js
Restart=on-failure
RestartSec=5
Environment=PORT=3001
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
